# templates/vault.service.erb
[Unit]
Description=vault

[Service]
ExecStart=/usr/bin/vault server -config=/etc/vault/settings.hcl
